# scripts/pump_pi_setup/systemd/sugar-uploader.service
[Unit]
Description=Sugar House Upload Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO_ROOT__
Environment=PYTHONUNBUFFERED=1
Environment=SUGAR_CONFIG_DIR=__REPO_ROOT__/config
ExecStart=__VENV_PATH__/bin/python __REPO_ROOT__/scripts/upload_service.py
Restart=always
RestartSec=2
StandardOutput=append:__LOG_PATH__
StandardError=append:__LOG_PATH__

[Install]
WantedBy=multi-user.target
